# scripts/com.hunter.update_charts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hunter.update_charts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/gianf/Desktop/Cursor/Hunter/scripts/update_charts.sh</string>
    </array>

    <!-- Run every hour at minute 0 (e.g., 1:00, 2:00, 3:00) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/gianf/Desktop/Cursor/Hunter</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/gianf/Desktop/Cursor/Hunter/.venv/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>MPLCONFIGDIR</key>
        <string>/Users/gianf/Desktop/Cursor/Hunter/.matplotlib</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gianf/Desktop/Cursor/Hunter/logs/update_charts.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gianf/Desktop/Cursor/Hunter/logs/update_charts.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
